# Youverse.Hosting/dotyoucore-watcher.path
[Path]
PathModified=/srv/sites/dotyoucore

[Install]
WantedBy=multi-user.target